# vault/vault.ini
[program:vault]
command=/bin/vault server -config=/etc/vault.conf
autostart=true
autorestart=false
stdout_logfile=/tmp/%(program_name)s.log
stderr_logfile=/tmp/%(program_name)s.log